# rhome_eod_webapp/ops/rhome_eod_webapp_ingest.service
[Unit]
Description=r:home Web EOD Dashboard (ingest)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/rhome_eod_webapp
EnvironmentFile=/etc/rhome_eod_webapp.env
ExecStart=/usr/bin/python3 -m app.ingest --mode live
